# src/conf/crontab
PATH=/usr/local/bin:/usr/bin:/usr/local/sbin:/usr/sbin
0 * * * * python3 /home/ec2-user/prod-batch/src/kick/UpdateContestUserKick.py
